# pyfluff.service
[Unit]
Description=PyFluff - Furby Connect Controller
Documentation=https://github.com/yourusername/PyFluff
After=network.target bluetooth.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/PyFluff
Environment="PATH=/home/pi/PyFluff/venv/bin"
# Uses default settings (0.0.0.0:8080). For custom settings, use:
# ExecStart=/home/pi/PyFluff/venv/bin/uvicorn pyfluff.server:app --host 0.0.0.0 --port 8080
ExecStart=/home/pi/PyFluff/venv/bin/python -m pyfluff.server
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
